# 990101/source/flask.service
[Unit]
Description = webrtc by flask

[Service]
Type=forking
ExecStart = /usr/bin/webrtc_Tony/flask_run.sh

[Install]
WantedBy = multi-user.target
